# misaka-core-v4/scripts/misaka-node.service
[Unit]
Description=MISAKA Network Node (Narwhal/Bullshark)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/misaka/misaka-public-node-linux-x86_64
ExecStart=/home/ubuntu/misaka/misaka-public-node-linux-x86_64/misaka-node \
    --config /home/ubuntu/misaka/misaka-public-node-linux-x86_64/config/public-node.toml \
    --data-dir /home/ubuntu/misaka/misaka-public-node-linux-x86_64/misaka-data \
    --genesis-path /home/ubuntu/misaka/misaka-public-node-linux-x86_64/config/genesis_committee.toml \
    --chain-id 2 \
    --advertise-addr 133.167.126.51:16110
EnvironmentFile=-/etc/misaka/node.env
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
LimitNOFILE=65535

StandardOutput=journal
StandardError=journal
SyslogIdentifier=misaka-node

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/misaka/misaka-public-node-linux-x86_64/misaka-data
PrivateTmp=true
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
