# relayer/misaka-relayer.service
[Unit]
Description=MISAKA Bridge Relayer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=misaka
Group=misaka
ExecStart=/usr/local/bin/misaka-relayer
Restart=always
RestartSec=10
EnvironmentFile=-/etc/misaka/relayer.env

# Local-only overrides; the rest of the runtime config comes from the env file.
Environment=RELAYER_KEYPAIR=/etc/misaka/relayer.json
Environment=PROCESSED_STORE=/var/lib/misaka/relayer-processed.json
Environment=RUST_LOG=info

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/misaka
PrivateTmp=true

[Install]
WantedBy=multi-user.target
